# supervisord.conf
[supervisord]
nodaemon=true
[program:nginx]
command=nginx -g "daemon off;"
[program:php-fpm]
command=docker-php-entrypoint php-fpm
